# supervisor.conf
[program:youkonger]
command=/root/youkonger/youkonger-bin
autostart=true
startsecs=1
startretries=3
stdout_logfile=/var/log/youkonger
stderr_logfile=/var/log/youkonger